# .config/systemd/user/trash-purge.service
[Unit]
Description=Purge trash files older than 30 days

[Service]
Type=oneshot
ExecStart=/usr/bin/trash-empty 30
